# supervisord.conf
[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root

[program:cqhttp]
directory=/usr/local/maidnana/go-cqhttp
command=/bin/sh -c "echo 02 | /usr/local/maidnana/go-cqhttp/go-cqhttp"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false

[program:nana]
directory=/usr/local/maidnana/MaidNanaGo
command=/usr/local/maidnana/MaidNanaGo/MaidNanaGo
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false